# infra/vm/systemd/rucksack-autopilot-v1-ZXJuaWVzZy9pZi1sZXR0ZXJzLWhvbWUtY291bGQtc2luZw-drain.timer
# Generated by rucksack. Runs in the trusted VM user's local timezone.
[Unit]
Description=Rucksack autopilot queue heartbeat
Documentation=man:systemd.timer(5)

[Timer]
OnBootSec=5m
OnUnitActiveSec=30m
Persistent=true
RandomizedDelaySec=5m
Unit=rucksack-autopilot-v1-ZXJuaWVzZy9pZi1sZXR0ZXJzLWhvbWUtY291bGQtc2luZw-drain.service

[Install]
WantedBy=timers.target
