# infra/vm/systemd/rucksack-codex-update.timer
# Generated by rucksack. Runs once daily with jitter and catches missed runs.
[Unit]
Description=Daily managed Codex and Claude Code CLI update timer for Rucksack

[Timer]
OnCalendar=daily
RandomizedDelaySec=2h
Persistent=true
Unit=rucksack-codex-update.service

[Install]
WantedBy=timers.target
