# infra/vm/systemd/rucksack-codex-update.service
# Generated by rucksack. Updates installed managed user-local agent CLIs.
[Unit]
Description=Daily managed Codex and Claude Code CLI update for Rucksack
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=6h
StartLimitBurst=4

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rucksack-codex-auto-update
Restart=on-failure
RestartSec=30m
RestartPreventExitStatus=20
TimeoutStartSec=15m
