# mac_bridge/launchd/com.erniesg.letters-home.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>KeepAlive</key>
    <true/>
    <key>Label</key>
    <string>com.erniesg.letters-home.bridge</string>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>-i</string>
        <string>HOME=/Users/letters-home</string>
        <string>PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <string>PYTHONPATH=/opt/letters-home</string>
        <string>/opt/python/bin/python3</string>
        <string>-m</string>
        <string>mac_bridge.server</string>
        <string>--repo-root</string>
        <string>/opt/letters-home</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/letters-home/.local/share/letters-home/bridge.stderr.log</string>
    <key>StandardOutPath</key>
    <string>/Users/letters-home/.local/share/letters-home/bridge.stdout.log</string>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>WorkingDirectory</key>
    <string>/opt/letters-home</string>
</dict>
</plist>
